# deploy/quadlet/moj-api.container
# moj-api.container — API MOJ (fcgiwrap + router.sh) como container rootless.
# Instale em ~/.config/containers/systemd/ e `systemctl --user daemon-reload`.
# O nginx do HOST (fora da imagem) serve web/ e faz fastcgi_pass ao socket unix em
#   /home/ribas/moj/run/fcgiwrap.sock  (= /data/run/fcgiwrap.sock dentro do container).
[Unit]
Description=MOJ API (fcgiwrap + router.sh)
After=network-online.target

[Container]
Image=localhost/moj-server:prod
# `podman auto-update` reinicia quando a imagem :prod local muda (build ou pull re-tagueado).
AutoUpdate=local
Exec=api
# Estado por volume, rótulo SELinux :z (SHARED — o nginx do host e o moj-judged mantêm acesso).
# NÃO usar :Z (privado) nem User=/UserNS=keep-id (quebraria o dono dos arquivos nos volumes).
Volume=/home/ribas/moj/run:/data/run:z
Volume=/home/ribas/moj/contests:/data/contests:z
Volume=/home/ribas/moj/moj-problems:/data/moj-problems:z
Volume=/home/ribas/moj/cdmoj/server/var/news:/data/news:z
Environment=INTAKE_MODE=queue

[Service]
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
