# deploy/quadlet/moj-judged.container
# moj-judged.container — daemon de julgamento (consome o spool, enfileira p/ o pull).
# Container SEPARADO do moj-api p/ restart independente (CLAUDE.md: "mexeu no judged -> reinicie").
# Compartilha os MESMOS volumes :z (o judged observa run/spool que a API escreve).
[Unit]
Description=MOJ judging daemon (queue/pull)
After=network-online.target

[Container]
Image=localhost/moj-server:prod
AutoUpdate=local
Exec=judged
Volume=/home/ribas/moj/run:/data/run:z
Volume=/home/ribas/moj/contests:/data/contests:z
Volume=/home/ribas/moj/moj-problems:/data/moj-problems:z
Volume=/home/ribas/moj/cdmoj/server/var/news:/data/news:z
Environment=INTAKE_MODE=queue
Environment=JUDGE_BACKEND=queue

[Service]
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
