# server/etc/systemd/moj-judged.service
[Unit]
Description=MOJ async judging daemon (consome o spool de submissões via inotify)
Documentation=file:///home/ribas/moj/docs/PLAN.md
After=network.target

[Service]
Type=simple
# Ajuste WorkingDirectory/MOJ_HOME por deploy.
Environment=MOJ_HOME=%h/moj
# Backend: mock (testes) | local (bubblewrap) | cluster (escalonador).
Environment=JUDGE_BACKEND=mock
# Override de caminhos, se quiser (senão usa server/etc/common.conf):
# Environment=CONTESTSDIR=%h/moj/contests
# Environment=RUNDIR=%h/moj/run
# Environment=JUDGE_MASTER=localhost:27000
# Environment=RESULT_SINK=localhost:28000
WorkingDirectory=%h/moj/server/daemons
ExecStart=/bin/bash %h/moj/server/daemons/judged.sh
Restart=on-failure
RestartSec=2
# Endurecimento básico (user service):
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
