# server/etc/systemd/moj-result-sink.service
[Unit]
Description=MOJ judge result sink (recebe veredictos por PUSH do worker/master)
Documentation=file:///home/ribas/moj/server/judge-gw/README.md
After=network.target

[Service]
Type=simple
Environment=RESULT_SINK_PORT=28000
# Environment=RUNDIR=%h/moj/run
WorkingDirectory=%h/moj/server/judge-gw
ExecStart=/bin/bash %h/moj/server/judge-gw/result-sink.sh
Restart=on-failure
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=default.target
