# server/etc/systemd/moj-agent@.service
[Unit]
# Agente de juiz MOJ (modelo PULL) — uma instância por CAPACIDADE: moj-agent@pos,
# moj-agent@gpu, moj-agent@cm, moj-agent@hu. Substitui root-daemon* + job-receiveitor*
# (sem porta de entrada, sem nc, sem poll-storm): só faz curl de SAÍDA p/ a API.
#   systemctl enable --now moj-agent@pos
Description=MOJ judge agent (pull) — capability %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Config por deploy: o agent.env (copiado de judge/etc/agent.env.sample) sobrescreve
# os defaults abaixo. O '-' torna o arquivo opcional.
EnvironmentFile=-/home/prof/judge/etc/agent.env
Environment=MOJ_API=http://localhost/api/v1
Environment=CAPABILITY=%i
Environment=WORKER_TOKEN_FILE=/home/prof/judge/etc/worker.token
Environment=MOJTOOLS_DIR=/home/prof/mojtools
# modelo cache: baixa o pacote por problema p/ o cache local (sem clonar repo/NFS)
Environment=JUDGE_CACHE=/home/prof/judge/cache/problems
Environment=HEARTBEAT_SECS=3
WorkingDirectory=/home/prof/judge/agent
ExecStart=/bin/bash /home/prof/judge/agent/moj-agent.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
