# server/etc/systemd/moj-gitea.service
[Unit]
Description=MOJ Gitea (store git da gestão de problemas — backend, sem acesso direto de autores)
Documentation=file:///home/ribas/moj/cdmoj/docs/DEPLOY-GITEA.md
After=network.target

[Service]
Type=simple
Environment=RUNDIR=%h/moj/run
WorkingDirectory=%h/moj/run/gitea
# O launcher é idempotente e faz exec do binário (vira o processo do serviço).
ExecStart=/usr/bin/env bash %h/moj/cdmoj/server/bin/start-gitea.sh
Restart=on-failure
RestartSec=3
# Endurecimento leve (user service): o Gitea só precisa do próprio home de trabalho.
NoNewPrivileges=true

[Install]
WantedBy=default.target
