# server/etc/systemd/moj-judged.service
[Unit]
Description=MOJ async judging daemon (consome o spool de submissões via inotify)
Documentation=file:///home/ribas/moj/cdmoj/docs/PLAN.md
After=network.target

[Service]
Type=simple
# Ajuste WorkingDirectory/MOJ_HOME por deploy.
Environment=MOJ_HOME=%h/moj/cdmoj
# Produção = PULL: o daemon enfileira e os juízes (moj-agent@) puxam o job.
Environment=INTAKE_MODE=queue
Environment=JUDGE_BACKEND=queue
# Dev/legado (intake inline, sem juízes): INTAKE_MODE=legacy + JUDGE_BACKEND=mock|local.
# Override de caminhos, se quiser (senão usa server/etc/common.conf):
# Environment=CONTESTSDIR=%h/moj/contests
# Environment=RUNDIR=%h/moj/run
WorkingDirectory=%h/moj/cdmoj/server/daemons
ExecStart=/bin/bash %h/moj/cdmoj/server/daemons/judged.sh
Restart=on-failure
RestartSec=2
# Endurecimento básico (user service):
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
