# server/etc/systemd/moj-bot.service
[Unit]
Description=MOJ Telegram bot (mojinho) — cliente da API, ENJAULADO com bwrap (SAMPLE)
Documentation=file:///home/ribas/moj/cdmoj/mojinho-bot/README.md
After=network.target
# Roda via run-caged.sh: jaula bwrap sem /home/workspace/contests/run/moj-problems.
# Segredos SÓ no dir vivo (~/mojinho-live: token + bot.conf) — nunca no repo/unit.

[Service]
Type=simple
ExecStart=/bin/bash %h/moj/cdmoj/mojinho-bot/run-caged.sh
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
